# severos.service
[Unit]
Description=SeverOS Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/workspace/SeverOS
ExecStart=/root/workspace/SeverOS/venv/bin/python /root/workspace/SeverOS/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
